# rootfs/etc/s6-overlay/s6-rc.d/logs/run
#!/command/execlineb -P
tail -n 0 -f /www/app/logs/php-error.log
