# rootfs/etc/s6-overlay/s6-rc.d/sshd/run
#!/command/with-contenv bash
[ "$RUNTIME_SSH" != true ] && exit 0
mkdir -p /run/sshd
/usr/sbin/sshd -D -e
