# rootfs/etc/s6-overlay/s6-rc.d/openlitespeed/run
#!/command/execlineb -P
/usr/local/lsws/bin/openlitespeed -d
